# systemd/wallbug.service
[Unit]
Description=Wall_Bug Daemon
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=wallbug
Group=wallbug
WorkingDirectory=/opt/wall_bug

Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/wall_bug
EnvironmentFile=-/etc/default/wallbug

ExecStart=/usr/bin/python3 -m wallbug.cli daemon --foreground
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
